# systemd/vigil-sweep.timer
[Unit]
Description=Run Vigil security sweep daily at 05:00

[Timer]
OnCalendar=*-*-* 05:00:00
Persistent=true
RandomizedDelaySec=30m

[Install]
WantedBy=timers.target
